# infra/systemd/poker-pg-retention.timer
[Unit]
Description=Prune Poker with Friends PostgreSQL backups daily

[Timer]
OnCalendar=*-*-* 03:30:00 UTC
RandomizedDelaySec=2m
Persistent=true
Unit=poker-pg-retention.service

[Install]
WantedBy=timers.target
